# quilt/contrib/supervisor/templates/program.conf
# This file is autogenerated. Manual changes may be overwritten!
[program:{{ name }}]
command={{ command }}
{% if numprocs > 1 -%}
process_name=%(group_name)_%(program_name)s_%(process_num)02d
{% endif -%}
numprocs={{ numprocs }}
priority={{ priority }}
autostart={{ autostart }}
autorestart={{ autorestart }}
startsecs={{ startsecs }}
startretries={{ retries }}
exitcodes={{ exitcodes }}
stopsignal={{ stopsignal }}
stopwaitsecs={{ stopwait }}
{% if program_user -%}
user={{ program_user }}
{% endif -%}
redirect_stderr={{ redirect_stderr }}
{% if stdout_logfile -%}
stdout_logfile={{ stdout_logfile }}
{% else -%}
stdout_logfile={{ log_dir }}/{{ name }}.out
{% endif -%}
stdout_logfile_maxbytes={{ stdout_logfile_maxsize }}
stdout_logfile_backups={{ stdout_logfile_keep }}
{% if stderr_logfile -%}
stderr_logfile={{ stderr_logfile }}
{% else -%}
stderr_logfile={{ log_dir }}/{{ name }}.err
{% endif -%}
stderr_logfile_maxbytes={{ stderr_logfile_maxsize }}
stderr_logfile_backups={{ stderr_logfile_keep }}
{% if environment -%}
environment={{ environment }}
{% endif -%}
{% if chdir -%}
directory={{ chdir }}
{% endif -%}
{% if umask -%}
umask={{ umask }}
{% endif -%}
